# home/amtea/.config/systemd/user/xdg-desktop-portal.service.d/override.conf
[Unit]
After=xdg-desktop-portal-hyprland.service
Requires=xdg-desktop-portal-hyprland.service
