# etc/systemd/system/getty@tty1.service.d/override.conf
[Unit]
Conflicts=poweroff.target reboot.target

[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin amtea --noclear %I $TERM
